# sys-apps/ryzen-stabilizator/files/ryzen-stabilizator
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/ryzen-stabilizator"

depend() {
	need localmount
	after modules-load
}

start() {
	ebegin "Initializing ryzen"
	${command} --config=/etc/ryzen-stabilizator/settings.toml >/dev/null 2>&1
	eend $?
}

stop() {
	# Nothing to stop
	return 0;
}
