# x11-misc/dotool/files/dotoold.service
[Unit]
Description=dotool daemon

[Service]
ExecStart=/usr/bin/dotoold
