# x11-misc/dotool/files/dotoold.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Starts dotoold service"
supervisor="supervise-daemon"
command="/usr/bin/dotoold"

reload() {
	ebegin "Reloading $RC_SVCNAME"
	/usr/bin/kill -HUP $MAINPID
	eend $?
}
