# www-apps/ttyd/files/ttyd.initd
#!/sbin/openrc-run
# shellcheck shell=sh
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/ttyd"
command_args="${TTYD_OPTS:--i lo -p 7681 -O} login"
command_background=1
pidfile="/run/ttyd.pid"

depend() {
	use net
}
