# app-admin/himitsu/files/himitsud.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
command=himitsud
