# www-apps/jackett-bin/files/jackett.initd-r1
#!/sbin/openrc-run
# Copyright 2020-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Jackett"
description="Jackett torrent tracker"
command="/opt/jackett/jackett"
command_background=true
command_user="jackett:jackett"
start_stop_daemon_args="--chdir /var/jackett --env HOME=/var/jackett"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
