# media-sound/mpd-discord-rpc/files/mpd-discord-rpc-1.11.0.user.init
#!/sbin/openrc-run

description="Discord Rich Presence for MPD"
supervisor="supervise-daemon"
command=/usr/bin/mpd-discord-rpc

depend() {
	use mpd
}
