# sys-apps/deepcool-digital-linux/files/deepcool-digital-linux.initd
#!/sbin/openrc-run

description="DeepCool Digital"
command="/usr/bin/deepcool-digital-linux"
command_args=""
command_background=1
pidfile="/run/deepcool-digital-linux.pid"
