# sys-apps/lact/files/lactd.init
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="AMDGPU Control Daemon"
command="/usr/bin/lact"
command_args="daemon"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
