# net-nntp/inn/files/ovsqlite.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

# shellcheck shell=sh

command=/opt/inn/bin/ovsqlite-server
pidfile=/run/news/ovsqlite.pid
