# app-mobilephone/pinephone-modem-scripts/files/pinephone-modem-scripts.pinephone-modem-setup.service
[Unit]
Description=Setup modem

[Service]
Type=simple
ExecStart=/usr/sbin/pinephone-modem-setup.sh

[Install]
WantedBy=graphical.target
