# sys-apps/plocate/files/updatedb.timer
[Unit]
Description=Update the plocate database daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
