# sys-power/disable-c6/files/disable-c6
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
command=/usr/sbin/zenstates
command_args="--c6-disable"
pidfile=
name="disable-c6"
description="Ryzen Disable C6"
